# systemd/server-bot.service
[Unit]
Description=Server management Telegram bot
After=network.target

[Service]
Type=simple
User=serverbot
Group=serverbot
ExecStart=/usr/bin/python3 /opt/bots/server/telegram-bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=server-bot
EnvironmentFile=/etc/bots/server.env

[Install]
WantedBy=multi-user.target
